# web2py/applications/ControleEstoque/cron/crontab
#crontab
* * * * * * root *applications/ControleEstoque/cron/validade.py
